# .devcontainer/autostart/filius.desktop
[Desktop Entry]
Type=Application
Name=Filius
Exec=filius
X-GNOME-Autostart-enabled=true